# Influxdb/supervisor.conf
[supervisord]
nodaemon=true

[unix_http_server]
file=/var/run//supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:influxdb]
priority=10
command=influxdb --config /opt/influxdb/shared/config.toml
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s-error.log
autostart=true
autorestart=true

[program:sensu-client]
priority=40
directory=/tmp
command=/opt/sensu/bin/sensu-client -c /etc/sensu/conf.d/influxdb-config.json -d /etc/sensu -e /etc/sensu/extensions -v -l /var/log/sensu/influxdb-client.log
user=root
autostart=true
autorestart=true
